# resolver.crontab
# The first element of the path is a directory where the debian-sa1
# script is located
PATH=/usr/lib/sysstat:/usr/sbin:/usr/sbin:/usr/bin:/sbin:/bin

# Update local DNS to support Kubernetes DNS queries
* * * * * root echo "nameserver $(crictl -r unix:///run/containerd/containerd.sock inspectp $(crictl -r unix:///run/containerd/containerd.sock pods --name coredns -q) | jq .status.network.ip | tr -d \")" > /etc/resolv.conf